# packages/hermes-image/rootfs/etc/profile.d/ghostship-internal.sh
if [ -r /run/ghostship/internal.env ]; then
  . /run/ghostship/internal.env
fi
