# packages/hermes-image/rootfs/etc/s6/services/ttyd/run
#!/bin/sh
set -eu

exec setpriv --reuid 1000 --regid 1000 --clear-groups --inh-caps -all /usr/local/bin/ghostship-hermes-runtime ttyd-service
