# packages/hermes-image/rootfs/etc/services.d/hermes-dashboard/run
#!/command/with-contenv bash
set -euo pipefail
export HOME="${HOME:-/home/hermes}"
export HERMES_HOME="${HERMES_HOME:-$HOME/.hermes}"
export GHOSTSHIP_DASHBOARD_HOST="${GHOSTSHIP_DASHBOARD_HOST:-127.0.0.1}"
export GHOSTSHIP_DASHBOARD_PORT="${GHOSTSHIP_DASHBOARD_PORT:-9119}"
exec /command/s6-setuidgid hermes /opt/hermes/venv/bin/hermes dashboard --host "$GHOSTSHIP_DASHBOARD_HOST" --port "$GHOSTSHIP_DASHBOARD_PORT" --no-open
