# packages/hermes-image/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bash
set -euo pipefail
exec nginx -g 'daemon off;'
