# packages/hermes-image/rootfs/etc/services.d/ghostship-router/run
#!/command/with-contenv bash
set -euo pipefail
export HOME="${HOME:-/home/hermes}"
export HERMES_HOME="${HERMES_HOME:-$HOME/.hermes}"
export GHOSTSHIP_ROUTER_HOST="127.0.0.1"
export GHOSTSHIP_ROUTER_PORT="8788"
if [ -r /run/ghostship/internal.env ]; then
  . /run/ghostship/internal.env
  export _GHOSTSHIP_ROUTER_API_KEY
fi
exec /command/s6-setuidgid hermes /opt/ghostship-router/venv/bin/ghostship-hermes-router --host "$GHOSTSHIP_ROUTER_HOST" --port "$GHOSTSHIP_ROUTER_PORT"
