# packages/hermes-image/workstation-seed/.config/systemd/user/ghostship-asset-refresh.timer
[Unit]
Description=Run ghostship-hermes asset refresh every 6 hours

[Timer]
OnBootSec=30m
OnUnitActiveSec=6h
RandomizedDelaySec=20m
Unit=ghostship-asset-refresh.service

[Install]
WantedBy=timers.target
