# packages/hermes-image/workstation-seed/.config/systemd/user/ghostship-asset-refresh.service
[Unit]
Description=ghostship-hermes mutable asset refresh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ghostship-hermes-runtime refresh-assets-once
