# packages/hermes-image/workstation-seed/.config/systemd/user/ghostship-opencode-model-refresh.service
[Unit]
Description=ghostship-hermes opencode model refresh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ghostship-hermes-runtime refresh-opencode-models-once
