# packages/hermes-image/workstation-seed/.config/systemd/user/ghostship-app-update.service
[Unit]
Description=ghostship-hermes agent app updater

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ghostship-hermes-runtime update-apps-once
