# packages/hermes-image/rootfs/etc/s6/services/profile-reconciler/run
#!/bin/sh
set -eu

exec /usr/local/bin/ghostship-hermes-runtime profile-reconciler-loop
